# ansible/kubernetes/roles/sladkoff.grafana-kiosk/templates/getty@tty1.service.j2
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin {{ grafana_kiosk_user }} --noclear %I 38400 linux
